# ft_services/srcs/influxdb/supervisord.conf
[supervisord]
nodaemon=true

[program:influxd]
command=/usr/bin/influxd
autorestart=true

[program:telegraf]
command=/usr/bin/telegraf
autorestart=true